# scheduling/com.concert-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.concert-scraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>concert-scraper</string>
        <string>scrape</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Code/concert-scraper</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/concert-scraper.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/concert-scraper.err</string>
</dict>
</plist>
